# aws-supervisord.conf
[supervisord]
nodaemon=true

[program:nginx_server]
user=root
command= /usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rqscheduler]
command=rqscheduler --host "%(ENV_REDIS_URL)s"
stopsignal=TERM
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app_server]
user=root
command=gunicorn --bind 0.0.0.0:5000 main:application
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asp_create_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" ASP_CREATE_REPORT --with-scheduler
autostart=true
autorestart=true

[program:asp_verify_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" ASP_VERIFY_REPORT --with-scheduler
autostart=true
autorestart=true

[program:item_master_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" ITEM_MASTER_REPORT --with-scheduler
autostart=true
autorestart=true

[program:order_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" ORDER_REPORT --with-scheduler
autostart=true
autorestart=true

# [program:settlement_report_v2_worker]
# user=root
# command=rq worker --url "redis://%(ENV_REDIS_URL)s" SETTLEMENT_REPORT_V2 --with-scheduler
# autostart=true
# autorestart=true

[program:sales_traffic_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" SALES_TRAFFIC_REPORT --with-scheduler
autostart=true
autorestart=true

[program:item_master_cogs_import_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" ITEM_MASTER_COGS_IMPORT --with-scheduler
autostart=true
autorestart=true

[program:export_csv_data_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" EXPORT_CSV --with-scheduler
autostart=true
autorestart=true

[program:finance_event_data_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" FINANCE_EVENT_LIST --with-scheduler
autostart=true
autorestart=true

[program:ledger_summary_data_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" LEDGER_SUMMARY_REPORT --with-scheduler
autostart=true
autorestart=true

[program:fba_returns_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" FBA_RETURNS_REPORT --with-scheduler
autostart=true
autorestart=true

[program:fba_reimbursements_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" FBA_REIMBURSEMENTS_REPORT --with-scheduler
autostart=true
autorestart=true

[program:fba_customer_shipment_sales_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" FBA_CUSTOMER_SHIPMENT_SALES_REPORT --with-scheduler
autostart=true
autorestart=true

[program:item_master_update_catalog_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" ITEM_MASTER_UPDATE_CATALOG --with-scheduler
autostart=true
autorestart=true

[program:az_ad_create_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" AZ_ADS_CREATE_REPORT --with-scheduler
autostart=true
autorestart=true

[program:az_ad_verify_report_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" AZ_ADS_VERIFY_REPORT --with-scheduler
autostart=true
autorestart=true

[program:az_asp_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" AZ_SPONSORED_PRODUCT --with-scheduler
autostart=true
autorestart=true

[program:az_asb_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" AZ_SPONSORED_BRAND --with-scheduler
autostart=true
autorestart=true

[program:az_aspd_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" AZ_SPONSORED_DISPLAY --with-scheduler
autostart=true
autorestart=true

[program:az_as_mr_pz_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" AZ_PERFORMANCE_ZONE --with-scheduler
autostart=true
autorestart=true

[program:email_worker]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" SES_EMAIL_DELIVERY --with-scheduler
autostart=true
autorestart=true

[program:subscription_check]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" SUBSCRIPTION_CHECK --with-scheduler
autostart=true
autorestart=true

[program:redis]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" REDIS --with-scheduler
autostart=true
autorestart=true

[program:sales_order_metrics]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" SALES_ORDER_METRICS --with-scheduler
autostart=true
autorestart=true

[program:fba_inventory]
user=root
command=rq worker --url "redis://%(ENV_REDIS_URL)s" FBA_INVENTORY --with-scheduler
autostart=true
autorestart=true
